# src/systemd/cmd.target
[Unit]
Description=Command Mode
Requires=sysinit.target cmd.service
After=sysinit.target cmd.service
AllowIsolate=yes
